# NIDS-Threat-Intel-Project/phase1-lab-setup/config/apt-update.service
[Unit]
Description=Periodic apt update for NIDS baseline traffic

[Service]
Type=oneshot
ExecStart=/usr/bin/apt update
