# NIDS-Threat-Intel-Project/phase1-lab-setup/config/ping-victim.timer
[Unit]
Description=Run ping burst every 20 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
